# wombott.ini
[program:wombott]
directory=%(ENV_HOME)s/wombott
command=%(ENV_HOME)s/wombott/venv/bin/python wombott.py
autostart=yes
autorestart=yes
startsecs=10
